# systemd/quant-scan-market-close.timer
[Unit]
Description=Schedule quant scan stop at 15:05 on weekdays

[Timer]
OnCalendar=Mon..Fri 15:05:00
Persistent=true
Unit=quant-scan-market-close.service

[Install]
WantedBy=timers.target
