# systemd/quant-scan-market-open.timer
[Unit]
Description=Schedule quant scan start at 09:20 on weekdays

[Timer]
OnCalendar=Mon..Fri 09:20:00
Persistent=true
Unit=quant-scan-market-open.service

[Install]
WantedBy=timers.target
