# systemd/quant-scan.service
[Unit]
Description=Guanlan AllTick Real-time Scanner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=quant
Group=quant
WorkingDirectory=/opt/guanlan-quant
Environment=HOME=/home/quant
Environment=PYTHONUNBUFFERED=1
ExecStartPre=+/usr/bin/install -d -o quant -g quant /opt/guanlan-quant/.guanlan/runtime
ExecStart=/bin/bash -lc 'echo $$ > /opt/guanlan-quant/.guanlan/runtime/scan.pid; exec /opt/guanlan-quant/start_alltick_multi_token_seconds.sh >> /opt/guanlan-quant/.guanlan/runtime/scan.log 2>&1'
ExecStopPost=/bin/bash -lc 'rm -f /opt/guanlan-quant/.guanlan/runtime/scan.pid'
Restart=always
RestartSec=3
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
